# etc/cron.d/bib-data-archive
0  */3 * * * root /usr/bin/systemctl restart bib-data-archive > /dev/null 2>&1
